# deploy_tools/gunicorn-django-showcase.haverford.edu.service
[Unit]
Description=Gunicorn server for django-showcase.haverford.edu

[Service]
Restart=on-failure
User=www-data
WorkingDirectory=/srv/sites/django-showcase.haverford.edu
EnvironmentFile=/srv/sites/django-showcase.haverford.edu/.env

ExecStart=/srv/sites/django-showcase.haverford.edu/.venv/bin/gunicorn \
    --bind unix:/tmp/django-showcase.haverford.edu.socket \
    showcase_site.wsgi:application

[Install]
WantedBy=multi-user.target
